# update-motd.d/20-help
#!/bin/bash
source /etc/os-release

printf "\n"
printf " * Documentation:  %s\n" "$HOME_URL"
printf " * Support:        %s\n" "$SUPPORT_URL"
printf "\n"

