# profile.d/00-pengwin.sh
#!/bin/sh
# bashsupport disable=BP5007

save_environment() {
  {
    echo "PATH='$PATH'"
    echo "WSL_DISTRO_NAME='$WSL_DISTRO_NAME'"
    echo "WSL_INTEROP='$WSL_INTEROP'"
    echo "WSL_SYSTEMD_EXECUTION_ARGS='$WSL_SYSTEMD_EXECUTION_ARGS'"
    echo "PULSE_SERVER='$PULSE_SERVER'"
    echo "WAYLAND_DISPLAY='$WAYLAND_DISPLAY'"
  } >"${systemd_saved_environment}"
}



setup_interop() {
  # shellcheck disable=SC2155,SC2012
  export WSL_INTEROP="$(ls -U /run/WSL/*_interop | tail -1)"
}



setup_display() {
  if [ -n "${XRDP_SESSION}" ]; then
    if [ -f "${systemd_saved_environment}" ]; then
      set -a
      # shellcheck disable=SC1090
      . "${systemd_saved_environment}"
      set +a
    fi

    if [ -n "${WSL_INTEROP}" ]; then
      export WSL2=1

      setup_interop
    fi

    unset WAYLAND_DISPLAY
    if [ -n "$SYSTEMD_PID" ]; then
      rm -f /run/user/"$(id -u)"/wayland* 2>/dev/null
    fi

    return
  fi

  if [ -n "${SSH_CONNECTION}" ]; then
    return
  fi

  # check whether it is WSL1 or WSL2
  if [ -n "${WSL_INTEROP}" ]; then
    #Export an environment variable for helping other processes
    export WSL2=1

    if [ -n "${DISPLAY}" ]; then

        uid="$(id -u)"

        user_path="/run/user/${uid}"
        if [ ! -d "${user_path}" ]; then
          sudo /usr/local/bin/create_userpath "${uid}" 2>/dev/null
        fi

        if [ -z "$SYSTEMD_PID" ]; then
          export XDG_RUNTIME_DIR="${user_path}"
        fi

        wslg_runtime_dir="/mnt/wslg/runtime-dir"

        ln -fs "${wslg_runtime_dir}"/wayland-0 "${user_path}"/ 2>/dev/null
        ln -fs "${wslg_runtime_dir}"/wayland-0.lock "${user_path}"/ 2>/dev/null

        pulse_path="${user_path}/pulse"
        wslg_pulse_dir="${wslg_runtime_dir}"/pulse

        if [ ! -d "${pulse_path}" ]; then
          mkdir -p "${pulse_path}" 2>/dev/null

          ln -fs "${wslg_pulse_dir}"/native "${pulse_path}"/ 2>/dev/null
          ln -fs "${wslg_pulse_dir}"/pid "${pulse_path}"/ 2>/dev/null

        elif [ -S "${pulse_path}/native" ]; then
          # Handle stale socket: remove it and recreate as symlink to WSLg pulse
          rm -f "${pulse_path}/native" 2>/dev/null
          ln -s "${wslg_pulse_dir}"/native "${pulse_path}"/ 2>/dev/null
        fi

        unset user_path
        unset wslg_runtime_dir
        unset wslg_pulse_dir
        unset pulse_path
        unset uid

      return
    fi

    # enable external x display for WSL 2
    route_exec=$(wslpath 'C:\Windows\system32\route.exe')

    if route_exec_path=$(command -v route.exe 2>/dev/null); then
      route_exec="${route_exec_path}"
    fi

    wsl2_d_tmp="$(eval "$route_exec print 2> /dev/null" | grep 0.0.0.0 | head -1 | awk '{print $4}')"

    if [ -n "${wsl2_d_tmp}" ]; then
      export DISPLAY="${wsl2_d_tmp}":0
    else
      wsl2_d_tmp="$(ip route | grep default | awk '{print $3; exit;}')"
      export DISPLAY="${wsl2_d_tmp}":0
    fi

    unset wsl2_d_tmp
    unset route_exec
  else
    # enable external x display for WSL 1
    export DISPLAY=localhost:0

    # Export an environment variable for helping other processes
    unset WSL2
  fi
}

setup_dbus() {
  # if dbus-launch is installed, then load it
  if ! (command -v dbus-launch >/dev/null); then
    return
  fi

  # Enabled via systemd
  if [ -n "${DBUS_SESSION_BUS_ADDRESS}" ]; then
    return
  fi

  dbus_pid="$(pidof -s dbus-daemon)"

  if [ -z "${dbus_pid}" ]; then
    dbus_env="$(timeout 2s dbus-launch --auto-syntax)"
    eval "${dbus_env}"

    dbus_env_file="/tmp/dbus_env_${DBUS_SESSION_BUS_PID}"
    echo "${dbus_env}" >"${dbus_env_file}"

    unset dbus_env
  else
    # Reuse existing dbus session
    dbus_env_file="/tmp/dbus_env_${dbus_pid}"
    if [ -f "${dbus_env_file}" ]; then
      eval "$(cat "${dbus_env_file}")"
    fi
  fi

  unset dbus_pid
  unset dbus_env_file
}

main() {
  # Only the default WSL user should run this script
  if ! (id -Gn | grep -c "adm.*sudo\|sudo.*adm" >/dev/null); then
    return
  fi

  systemd_saved_environment="$HOME/.systemd.env"

  SYSTEMD_PID="$(ps -C systemd -o pid= | head -n1)"

  setup_display

  if [ -z "$SYSTEMD_PID" ]; then
    setup_dbus
  fi

  # speed up some GUI apps like gedit
  export NO_AT_BRIDGE=1

  # Fix 'clear' scrolling issues
  alias clear='clear -x'

  # Custom aliases
  alias ll='ls -al'
  alias winget='powershell.exe winget'
  alias wsl='wsl.exe'

  if [ -n "${WSL2}" ]; then
    #Setup video acceleration
    export VDPAU_DRIVER=d3d12
    export LIBVA_DRIVER_NAME=d3d12
    sudo /usr/local/bin/pengwin-load-vgem-module

    # Setup Gallium Direct3D 12 driver
    export GALLIUM_DRIVER=d3d12
  fi

  if [ -z "$SYSTEMD_PID" ]; then

    save_environment

  elif [ -n "$SYSTEMD_PID" ] && [ "$SYSTEMD_PID" -eq 1 ] && [ -f "$HOME/.systemd.env" ] && [ -n "$WSL_SYSTEMD_EXECUTION_ARGS" ]; then
    # Only if built-in systemd was started
    set -a
    # shellcheck disable=SC1090
    . "${systemd_saved_environment}"
    set +a

    setup_interop
  fi

  # Check if we have Windows Path
  if [ -z "$WIN_HOME" ] && (command -v cmd.exe >/dev/null 2>&1); then

    # Create a symbolic link to the windows home

    # Here has an issue: %HOMEDRIVE% might be using a custom set location
    # moving cmd to where Windows is installed might help: %SYSTEMDRIVE%
    wHomeWinPath=$(cmd-exe /c 'cd %SYSTEMDRIVE%\ && echo %HOMEDRIVE%%HOMEPATH%' 2>/dev/null | tr -d '\r')

    if [ ${#wHomeWinPath} -le 3 ]; then #wHomeWinPath contains something like H:\
      wHomeWinPath=$(cmd-exe /c 'cd %SYSTEMDRIVE%\ && echo %USERPROFILE%' 2>/dev/null | tr -d '\r')
    fi

    # shellcheck disable=SC2155
    export WIN_HOME="$(wslpath -u "${wHomeWinPath}")"

    win_home_lnk=${HOME}/winhome
    if [ ! -e "${win_home_lnk}" ]; then
      ln -s -f "${WIN_HOME}" "${win_home_lnk}" >/dev/null 2>&1
    fi

    unset win_home_lnk

  fi
}

main "$@"
