# debian/pengwin-base/etc/update-motd.d/00-header
#!/bin/sh

. /etc/os-release

printf "\n"
printf "Welcome to %s (%s), Based on %s %s\n" "${PRETTY_NAME}" "${PENGWIN_VERSION}" "${ID}" "${VERSION}"
